# deploy/firepulse.service
[Unit]
Description=FirePulse Heart Rate Monitoring Server
After=network.target

[Service]
Type=simple
User=firepulse
Group=firepulse
WorkingDirectory=/opt/firepulse/packages/backend
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=HTTP_PORT=3001
Environment=UDP_PORT=41234
Environment=DATABASE_URL=file:/opt/firepulse/data/firepulse.db
Environment=JWT_SECRET=firepulse-appliance-local

# Logging (stdout/stderr -> journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=firepulse

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/firepulse/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
